# drivers/r-01/upico.service
[Unit]
Description=upico

[Service]
Type=simple
ExecStart=/usr/local/bin/upico service r01

[Install]
WantedBy=multi-user.target